# root/etc/s6-overlay/s6-rc.d/adguardhome-sync/run
#!/command/execlineb -P
with-contenv
s6-setuidgid bsd
/usr/local/bin/adguardhome-sync run --config /config/adguardhome-sync.yaml
